# packages/keyszer/keyszer-0.6.0.tar.gz/keyszer-0.6.0/example/keyszer.service
[Unit]
Description=keyszer

[Service]
Type=simple
KillMode=process
ExecStart=keyszer --quiet --watch
Restart=on-failure
RestartSec=3
Environment=DISPLAY=:0
User=keymapper
Group=input

[Install]
WantedBy=graphical.target

# setfacl -m user:keymapper:rw /dev/uinput
# xhost +SI:localuser:keymapper
# KERNEL=="uinput", SUBSYSTEM=="misc", RUN+="/usr/bin/setfacl -m user:keymapper:rw /dev/uinput"
